# assets/etc/supervisor.d/services.ini
[supervisord]
user = root

[program:dnsmasq]
command     = /usr/sbin/dnsmasq -d
autostart   = true
autorestart = true

[program:tinyproxy]
command = /usr/sbin/tinyproxy -d
autostart   = true
autorestart = true

